# corvus.service
[Unit]
Description=Corvus VM Management Daemon
After=network.target

[Service]
Type=exec
# Uses Unix socket by default: $XDG_RUNTIME_DIR/corvus/corvus.sock
ExecStart=%h/.local/bin/corvus \
    --database postgresql://localhost/corvus

Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
