# systemd/corvus-netd.service
[Unit]
Description=Corvus privileged network agent
After=network-pre.target
Before=network.target
Wants=network-pre.target

[Service]
Type=exec
ExecStart=/usr/local/bin/corvus-netd --log-level info
Restart=on-failure
RestartSec=5

# Phase 1 stub: the agent doesn't yet touch the kernel, so a stricter
# sandbox would be safe. Phase 2 needs CAP_NET_ADMIN + raw netlink +
# /dev/net/tun + nft + dnsmasq spawning, so we leave the sandbox
# permissive for now and tighten in Phase 2.

[Install]
WantedBy=multi-user.target
