# yaml/corvus-test-node/systemd/corvus-test.service
[Unit]
Description=Corvus daemon (integration-test target)
# Unprivileged daemon. All privileged network operations are
# delegated to corvus-netd.service.
# The binary lives on virtiofs (mounted at /opt/corvus/bin by
# opt-corvus-bin.mount). Postgres must be reachable.
# `RequiresMountsFor` makes systemd auto-pull the mount unit AND
# fail this service cleanly with a clear journal message if the
# mount itself failed — better than the old `ConditionPathExists`
# silently skipping the service.
After=network-online.target postgresql-18.service opt-corvus-bin.mount corvus-netd.service
Wants=network-online.target postgresql-18.service corvus-netd.service
Requires=opt-corvus-bin.mount
RequiresMountsFor=/opt/corvus/bin

[Service]
Type=simple
User=corvus
Group=corvus
SupplementaryGroups=kvm
RuntimeDirectory=corvus
Environment=XDG_RUNTIME_DIR=/run/corvus
ExecStart=/opt/corvus/bin/corvus \
    --tcp \
    --host 0.0.0.0 \
    --port 9876 \
    --database postgresql://corvus@/corvus \
    --log-level info
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
