# yaml/corvus-test-node/systemd/corvus-tcp-relay.service
[Unit]
Description=VSOCK-to-TCP relay for the inner Corvus daemon
# No dependency on corvus-test.service: the harness deploys the
# daemon's mTLS cert over SSH-over-VSOCK at class fixture time
# and only THEN enables corvus-test. A transitive `Wants=` on
# corvus-test here would crash-loop the daemon at every test
# image boot (no cert in /etc/corvus yet) until StartLimit
# kicks in, and on agents-only nodes it would never recover.
# The relay tolerates an initially-absent TCP listener (socat
# retries on the FORK-side connect failure).

[Service]
Type=simple
ExecStart=/usr/bin/socat \
    VSOCK-LISTEN:9876,reuseaddr,fork \
    TCP:127.0.0.1:9876
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target
