# yaml/corvus-test-node/systemd/corvus-nodeagent.service
[Unit]
Description=Corvus per-host node agent
# Same virtiofs mount as the daemon; no Postgres dependency.
After=opt-corvus-bin.mount
Requires=opt-corvus-bin.mount
RequiresMountsFor=/opt/corvus/bin

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/corvus/bin/corvus-nodeagent \
    --host 127.0.0.1 \
    --port 9878 \
    --log-level info
KillSignal=SIGTERM
TimeoutStopSec=30s
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
