# python/corvus_admin/_systemd_templates/corvus-nodeagent.service.j2
[Unit]
Description=Corvus per-host node agent
After=network-pre.target
Before=network.target
Wants=network-pre.target

[Service]
Type=exec
{%- if install_mode == "system" %}
# Match the daemon's runtime directory. The agent computes per-VM
# QEMU socket paths from XDG_RUNTIME_DIR; without this the daemon's
# surviving direct-QMP calls (disk hotplug, build key injection,
# send-ctrl-alt-del) miss the sockets because the two processes
# resolve different paths.
RuntimeDirectory=corvus
Environment=XDG_RUNTIME_DIR=/run/corvus
{%- endif %}
ExecStart={{ binary_path }} --log-level {{ log_level }}
Restart=on-failure
RestartSec=5

[Install]
{%- if install_mode == "user" %}
WantedBy=default.target
{%- else %}
WantedBy=multi-user.target
{%- endif %}
